# opencode/rootfs/etc/services.d/opencode/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant OpenCode - Finish Script
# ==============================================================================

bashio::log.info "Stopping OpenCode..."
