# it-tools/rootfs/etc/services.d/it-tools/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant IT Tools - Finish Script
# ==============================================================================

bashio::log.info "Stopping IT Tools..."
