# openclaw/rootfs/etc/services.d/openclaw/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant OpenClaw - Finish Script
# ==============================================================================

bashio::log.info "Stopping OpenClaw..."
