# ollama-container/rootfs/etc/services.d/ollama-container/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Ollama Container - Finish Script
# ==============================================================================

bashio::log.info "Stopping Ollama Container..."
