# changedetection/rootfs/etc/services.d/changedetection/run
#!/usr/bin/with-contenv bashio

declare base_url
declare notification_url
declare fetch_workers
declare data_dir
declare log_level

base_url=$(bashio::config 'base_url')
notification_url=$(bashio::config 'notification_url')
fetch_workers=$(bashio::config 'fetch_workers' '2')
data_dir=$(bashio::config 'data_dir' '/share/changedetection')
log_level=$(bashio::config 'log_level' 'info')

export PORT=5000
export DATASTORE_PATH="${data_dir}"
export FETCH_WORKERS="${fetch_workers}"

if bashio::var.has_value "${base_url}"; then
    export BASE_URL="${base_url}"
fi
if bashio::var.has_value "${notification_url}"; then
    export NOTIFICATION_URL="${notification_url}"
fi

mkdir -p "${data_dir}"

bashio::log.info "Starting Changedetection.io..."
bashio::log.info "Data Dir: ${data_dir}"
bashio::log.info "Fetch Workers: ${fetch_workers}"

cd /app || exit 1
exec python3 changedetection.py
