# grocy/rootfs/etc/services.d/grocy/run
#!/usr/bin/with-contenv bashio

declare timezone
declare log_level

timezone=$(bashio::config 'timezone' 'Etc/UTC')
log_level=$(bashio::config 'log_level' 'info')

export TZ="${timezone}"

mkdir -p /run/nginx /var/www/html/data
cat > /etc/nginx/http.d/default.conf <<'EOF'
server {
    listen 9283;
    root /var/www/html/public;
    index index.php index.html;

    location / {
        try_files $uri $uri/ /index.php?$query_string;
    }

    location ~ \.php$ {
        fastcgi_pass 127.0.0.1:9000;
        fastcgi_param SCRIPT_FILENAME $document_root$fastcgi_script_name;
        include fastcgi_params;
    }
}
EOF

bashio::log.info "Starting Grocy..."
bashio::log.info "Timezone: ${timezone}"

php-fpm83 &
exec nginx -g "daemon off;"
